# contrib/freebsd/files/cloudflare_ddns.in
#!/bin/sh
#
# PROVIDE: cloudflare_ddns
# REQUIRE: DAEMON NETWORKING
# KEYWORD: shutdown
#
# Add the following to /etc/rc.conf[.local] to enable this service:
#
#   cloudflare_ddns_enable="YES"
#
# Optional tunables (with defaults):
#
#   cloudflare_ddns_user="cloudflare_ddns"
#   cloudflare_ddns_config="%%PREFIX%%/etc/cloudflare-register/config.json"
#   cloudflare_ddns_command="%%PREFIX%%/bin/cloudflare-register"
#
# The config file is JSON; the service reads the path from the
# CLOUDFLARE_REGISTER_CONFIG environment variable. Generate it with
# 'cloudflare-register init --system' (then move/point it as needed).

. /etc/rc.subr

name="cloudflare_ddns"
rcvar="${name}_enable"

load_rc_config "$name"

: "${cloudflare_ddns_enable:="NO"}"
: "${cloudflare_ddns_user:="cloudflare_ddns"}"
: "${cloudflare_ddns_config:="%%PREFIX%%/etc/cloudflare-register/config.json"}"
: "${cloudflare_ddns_command:="%%PREFIX%%/bin/cloudflare-register"}"

pidfile="/var/run/${name}.pid"
required_files="${cloudflare_ddns_config}"

cloudflare_ddns_env="CLOUDFLARE_REGISTER_CONFIG=${cloudflare_ddns_config}"

# Run the foreground console script under daemon(8), which supervises it,
# drops privileges, and writes the pidfile.
command="/usr/sbin/daemon"
command_args="-f -p ${pidfile} -u ${cloudflare_ddns_user} ${cloudflare_ddns_command} service"

run_rc_command "$1"
