# contrib/freebsd/pkg-plist
@sample etc/cloudflare-register/cloudflare-register.env.sample
bin/cloudflare-register
@sample etc/rc.d/cloudflare_ddns
lib/python3.11/site-packages/cloudflare_register/__init__.py
lib/python3.11/site-packages/cloudflare_register/cli.py
lib/python3.11/site-packages/cloudflare_register/config.py
lib/python3.11/site-packages/cloudflare_register/exceptions.py
lib/python3.11/site-packages/cloudflare_register/ip_detection.py
lib/python3.11/site-packages/cloudflare_register/logging_setup.py
lib/python3.11/site-packages/cloudflare_register/persistence.py
lib/python3.11/site-packages/cloudflare_register/sync.py
lib/python3.11/site-packages/cloudflare_register/__main__.py
lib/python3.11/site-packages/cloudflare_register/providers/__init__.py
lib/python3.11/site-packages/cloudflare_register/providers/base.py
lib/python3.11/site-packages/cloudflare_register/providers/cloudflare.py
lib/python3.11/site-packages/cloudflare_register/providers/factory.py
lib/python3.11/site-packages/cloudflare_register/web/__init__.py
lib/python3.11/site-packages/cloudflare_register/web/app.py
lib/python3.11/site-packages/cloudflare_register/web/csrf.py
lib/python3.11/site-packages/cloudflare_register/tui/__init__.py
lib/python3.11/site-packages/cloudflare_register/tui/app.py
share/licenses/cloudflare-register/LICENSE
share/licenses/cloudflare-register/catalog.mk
@comment generated with bmake plist-mode
