# debian/cloudflare-register.cloudflare-ddns.service
[Unit]
Description=cloudflare-register (smart Cloudflare DynDNS)
Documentation=https://github.com/cloudbsdorg/cloudflare-register
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cloudflare-ddns
Group=cloudflare-ddns
WorkingDirectory=/var/lib/cloudflare-register
# cloudflare-register reads its JSON config from the path in
# CLOUDFLARE_REGISTER_CONFIG (default resolution would find
# /etc/cloudflare-register.json anyway; set it explicitly so operators can
# override it with a drop-in without editing this unit).
Environment=CLOUDFLARE_REGISTER_CONFIG=/etc/cloudflare-register.json
ExecStart=/usr/bin/cloudflare-register service
Restart=on-failure
RestartSec=10

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/cloudflare-register /var/log/cloudflare-register
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
SystemCallArchitectures=native
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
SystemCallFilter=@system-service
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
RestrictSUIDSGID=true
RemoveIPC=true
ProtectProc=invisible
UMask=0077

[Install]
WantedBy=multi-user.target
